# etc/systemd/system/mosh_packets.service
[Unit]
Description=set environment variable for mosh packets received

[Service]
Type=simple
ExecStart=/home/trey/bin/mosh_packets.py
